# bin/routechoices-spot.service
[Unit]
Description="routechoices spot crawler"
Wants=network-online.target
After=network-online.target

[Service]
Type=Simple
ExecStart=/apps/routechoices-server/env/bin/python /apps/routechoices-server/manage.py run_spot_crawler
Restart=always
Environment="PATH=/apps/routechoices-server/env/bin/"
WorkingDirectory=/apps/routechoices-server/
KillMode=mixed
TimeoutSec=5
PrivateTmp=true

[Install]
WantedBy=default.target
