# bin/routechoices-tornado.service
[Unit]
Description="routechoices tornado server"
Wants=network-online.target
After=network-online.target

[Service]
Type=Simple
ExecStart=parallel -j 4 -n0 /apps/routechoices-server/bin/run_tornado.sh ::: {1..4}
Restart=always
Environment="PATH=/apps/routechoices-server/env/bin/"
WorkingDirectory=/apps/routechoices-server/
KillMode=mixed
TimeoutSec=5
PrivateTmp=true

[Install]
WantedBy=default.target
